# deploy/systemd/caimeo-backend.service
[Unit]
Description=CAIMEO FastAPI backend
After=network-online.target
Wants=network-online.target

[Service]
User=pi
WorkingDirectory=/home/pi/CaimeoV1
EnvironmentFile=/home/pi/CaimeoV1/.env
ExecStart=/home/pi/CaimeoV1/venv/bin/uvicorn server:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
